# etc/xdg/autostart/calamares.desktop
[Desktop Entry]
Type=Application
Name=Calamares Installer
Exec=calamares
Icon=calamares
Terminal=false
StartupNotify=false
X-KDE-autostart-after=panel
